# holer-server/src/main/resources/bin/holer.service
[Unit]
Description=Holer
After=network.target

[Service]
Type=forking
ExecStart=@HOLER_HOME@/bin/holer start
ExecStop=@HOLER_HOME@/bin/holer stop

[Install]
WantedBy=multi-user.target
